# crontabs/www-data
* * * * * cd /var/www/toran && php bin/cron -v
